# pdns/supervisord.conf
[supervisord]
nodaemon=true

[program:pdns]
command=/usr/sbin/pdns_server --guardian=yes
